# supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
command=/opt/elasticsearch/bin/elasticsearch -Des.insecure.allow.root=true

[program:kibana]
command=/opt/kibana/bin/kibana
